# MAlab3/cron/cron/crontab
# Запуск скрипта каждые 5 часов
*/5 * * * * root /bin/bash /scripts/backup_script.sh
